# supervisor/cookies-task.conf
[program:cookies-task]
directory=/apps/cookies/
command = python /apps/cookies/task.py
autostart=true
autorestart = true
