# supervisor/cookies-proxy.conf
[program:cookies-proxy]
directory=/apps/cookies
command = python /apps/cookies/proxy.py
autostart=true
autorestart = true
